# files/etc/hotplug.d/iface/90-ipv6-pd-monitor
#!/bin/sh

WAN_IFACE="wan6"
RECOVERY_DELAY=3
INITIAL_WAIT=30
LOG_FILE="/var/log/ipv6-pd-monitor.log"
PD_CACHE="/tmp/ipv6_pd_last_prefix"

log() {
    logger -t ipv6-pd-monitor "$1"
    echo "$(date '+%Y-%m-%d %H:%M:%S') [Hotplug] $1" >> "$LOG_FILE"
}

# Hotplug trigger condition check
[ "$INTERFACE" != "$WAN_IFACE" ] && exit 0

case "$ACTION" in
    ifup|ifupdate)
        STATUS_DATA=$(ifstatus $WAN_IFACE)
        IS_UP=$(echo "$STATUS_DATA" | jsonfilter -e '@["up"]')
        if [ "$IS_UP" = "true" ]; then
            sleep $INITIAL_WAIT
            CURRENT_PD=$(ifstatus $WAN_IFACE | jsonfilter -e '@["ipv6-prefix"][*].address' 2>/dev/null | head -n 1)
            if [ -z "$CURRENT_PD" ]; then
                log "🔥 Warning: $WAN_IFACE is UP, but no IPv6-PD prefix found. Restarting interface..."
                ifdown $WAN_IFACE
                sleep $RECOVERY_DELAY
                ifup $WAN_IFACE
            else
                LAST_PD=$(cat "$PD_CACHE" 2>/dev/null)
                if [ "$CURRENT_PD" != "$LAST_PD" ]; then
                    log "💡 IPv6-PD Change Detected: [${LAST_PD:-INIT}] -> [$CURRENT_PD]"
                    echo "$CURRENT_PD" > "$PD_CACHE"
                    if [ -x "/etc/init.d/bypass_guard" ]; then
                        log "🔄 PD changed, refreshing bypass_guard rules..."
                        /etc/init.d/bypass_guard refresh_nft
                    fi
                fi
            fi
        fi
        ;;
esac

exit 0
