# files/etc/hotplug.d/iface/90-ipv6-pd-monitor
#!/bin/sh

WAN_IFACE="wan6" 
RECOVERY_DELAY=3 
INITIAL_WAIT=5  
LOG_FILE="/var/log/ipv6-pd-monitor.log"

log() {
    logger -t ipv6-pd-monitor "$1"
    echo "$(date '+%Y-%m-%d %H:%M:%S') [Hotplug] $1" >> "$LOG_FILE"
}

# Hotplug trigger condition check
[ "$INTERFACE" != "$WAN_IFACE" ] && exit 0

case "$ACTION" in
    ifup|ifupdate)
        # Check only when the interface logic is reported up.
        if (ifstatus $WAN_IFACE | grep -q '"up": true'); then
            log "$WAN_IFACE STATUS ($ACTION). Wait for $INITIAL_WAIT seconds to get PD for DHCPv6 client..."
            
            sleep $INITIAL_WAIT

            # --- Key check logic: directly check the PD prefix obtained by the wan6 interface ---
            PD_PREFIX=$(ifstatus $WAN_IFACE | jsonfilter -e '@["ipv6-prefix"][*].address' 2>/dev/null | head -n 1)
            
            if [ -z "$PD_PREFIX" ]; then
                log "🔥 Warning: $WAN_IFACE status UP, but the IPv6-PD prefix is not detected."
                log "Try to restart the $WAN_IFACE interface for recovery..."
                
                ifdown $WAN_IFACE
                sleep $RECOVERY_DELAY
                ifup $WAN_IFACE
                
                log "Restarting $WAN_IFACE. Wait for netifd to re-obtain PD..."
            else
                log "✅ IPv6-PD is obtained: $PD_PREFIX. No need to operate."
            fi
        fi
        ;;
esac

exit 0
