# files/etc/hotplug.d/iface/90-ipv6-pd-monitor
#!/bin/sh
# 文件名: /etc/hotplug.d/iface/90-ipv6-pd-monitor
# 功能: 监听 wan6 状态变化，在 PD 缺失时自动重启 wan6 接口进行恢复。
# 检查逻辑: 直接检查 wan6 接口的 ipv6-prefix 字段是否为空。

# --- 变量定义 ---
WAN_IFACE="wan6" 
RECOVERY_DELAY=3 
INITIAL_WAIT=5  
LOG_FILE="/var/log/ipv6-pd-monitor.log"

log() {
    logger -t ipv6-pd-monitor "$1"
    echo "$(date '+%Y-%m-%d %H:%M:%S') [Hotplug] $1" >> "$LOG_FILE"
}

# 1. Hotplug 触发条件检查
[ "$INTERFACE" != "$WAN_IFACE" ] && exit 0

case "$ACTION" in
    ifup|ifupdate)
        # 仅在接口逻辑上报 Up 时才检查
        if (ifstatus $WAN_IFACE | grep -q '"up": true'); then
            log "检测到 $WAN_IFACE 接口成功连接 ($ACTION)。等待 $INITIAL_WAIT 秒给 DHCPv6 客户端获取 PD..."
            
            sleep $INITIAL_WAIT

            # --- 关键检查逻辑：直接检查 wan6 接口获取的 PD 前缀 ---
            # 提取 wan6 接口上所有 'ipv6-prefix' 字段的地址。
            # 如果成功获取 PD，此数组非空；如果失败，则为空。
            PD_PREFIX=$(ifstatus $WAN_IFACE | jsonfilter -e '@["ipv6-prefix"][*].address' 2>/dev/null | head -n 1)
            
            # --- 3. 判断是否需要恢复 ---
            if [ -z "$PD_PREFIX" ]; then
                log "🔥 警告：$WAN_IFACE 状态 UP，但未检测到 IPv6-PD 前缀。"
                log "尝试重启 $WAN_IFACE 接口进行恢复..."
                
                # 执行 ifdown/ifup wan6 恢复操作
                ifdown $WAN_IFACE
                sleep $RECOVERY_DELAY
                ifup $WAN_IFACE
                
                log "已发送重启 $WAN_IFACE 命令。等待 netifd 重新获取 PD..."
            else
                log "✅ 检查正常：IPv6-PD 正常。已获取前缀: $PD_PREFIX/64。无需操作。"
            fi
        fi
        ;;
esac

exit 0
